# macos/launchd/com.pager.agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pager.agent</string>

    <!--
      Combined session + watchdog agent. One Login Items entry.

      `pager watchdog claude` already handles "session missing -> spawn it"
      via its existing restart path, so a single periodic agent covers
      both initial spawn (RunAtLoad fires immediately on load/login) and
      ongoing health checks (StartInterval every 70 s).
    -->

    <!--
      AssociatedBundleIdentifiers (macOS Ventura 13+) tells the Login
      Items UI which bundle to fetch the icon + display name from. The
      bundle ID matches CFBundleIdentifier in macos/pager.app/Contents/
      Info.plist. Without this key, Login Items shows the generic "Item
      from unidentified developer" + exec icon regardless of what the
      ProgramArguments path points at.
    -->
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.pager.agent</string>
    </array>

    <!--
      ProgramArguments points at the .app's launcher (a thin shell shim
      that exec's bin/pager). The Login Items icon resolution is handled
      separately via AssociatedBundleIdentifiers above.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>__PAGER_ROOT__/macos/pager.app/Contents/MacOS/pager</string>
        <string>watchdog</string>
        <string>claude</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>70</integer>

    <key>StandardOutPath</key>
    <string>__PAGER_ROOT__/logs/launchd-agent.out</string>
    <key>StandardErrorPath</key>
    <string>__PAGER_ROOT__/logs/launchd-agent.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__BREW_BIN__:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__USER_HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__USER_HOME__</string>
        <key>PAGER_ROOT</key>
        <string>__PAGER_ROOT__</string>
    </dict>
</dict>
</plist>
