# linux/systemd/pager-watch.service
[Unit]
Description=pager: watchdog tick — log session state and restart if dead
Documentation=https://github.com/jawwadzafar/pager
After=pager.service

[Service]
Type=oneshot
ExecStart=__PAGER_ROOT__/bin/pager watchdog default-boot
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PAGER_ROOT=__PAGER_ROOT__
