# linux/systemd/pager.service
[Unit]
Description=pager: persistent Claude Code session over Remote Control
Documentation=https://github.com/jawwadzafar/pager
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=__PAGER_ROOT__/bin/pager start default-boot
ExecStop=__PAGER_ROOT__/bin/pager stop default-boot
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=PAGER_ROOT=__PAGER_ROOT__

[Install]
WantedBy=default.target
